# oom-logger.timer
[Unit]
Description=Run OOM logger every second

[Timer]
OnBootSec=1
OnUnitActiveSec=1
AccuracySec=100ms

[Install]
WantedBy=timers.target
